# support/systemd.service
[Unit]
Description=Middleware for handling BLE button on Dell stylus

[Service]
Type=simple
ExecStart=/usr/local/bin/bluetooth-buttons
Nice=5
